# battalert.cron
# /etc/crond.d/battalert: crontab entries for Battery Alert

SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# m h dom mon dow user  command
  * * *   *   *   root  /usr/bin/battalert
